# vendor/lowkey/bootstraps/optional/idle-shutdown/systemd/idle-check.timer
[Unit]
Description=Agent idle check timer — every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=10s
Persistent=true

[Install]
WantedBy=timers.target
